# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/init-cloudflared-log/up
/bin/bash -c /etc/s6-overlay/scripts/00-cloudflared-log.sh