# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Cloudflared
#
# Configures the nginx service as workaround to get live streaming logs working
# see https://github.com/brenner-tobias/addon-cloudflared/discussions/744
# ==============================================================================
# set up and run nginx
# ==============================================================================

if bashio::config.false 'use_builtin_proxy'; then
    bashio::log.info "Using Cloudflared without the built-in NGINX proxy"
    # Send a signal to finish script to avoid halting container if the built-in proxy was disabled
    touch /dev/shm/no_built_in_proxy
    # Tell S6-Overlay not to restart this service
    s6-svc -O .
    return "${__BASHIO_EXIT_OK}"
fi

bashio::log.debug "Merging options & variables for template"
JSON_CONF=$(jq -n \
    --arg port "$(bashio::core.port)" \
    --arg ssl "$(bashio::core.ssl)" \
    '{port: $port, ssl: ($ssl | fromjson)}')
bashio::log.debug "Generating nginx.conf from template in /etc/nginx/servers/homeassistant.conf"
bashio::log.debug "JSON_CONF: ${JSON_CONF}"
tempio \
    -template /etc/nginx/template/homeassistant.gtpl \
    -out /etc/nginx/servers/homeassistant.conf \
    <<<"${JSON_CONF}"
bashio::log.debug "Generated nginx.conf:"
bashio::log.debug "$(cat /etc/nginx/servers/homeassistant.conf)"

# start server
bashio::log.info "Running nginx..."
exec nginx -c /etc/nginx/nginx.conf
