# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Connects and runs the Cloudfalred tunnel for HomeAssistant
# ==============================================================================
certificate="/data/cert.pem"
config="/data/config.yml"
tunnel_name="$(bashio::config 'tunnel_name')"

bashio::log.info "Connecting Cloudflared Tunnel..."

exec cloudflared --origincert=${certificate} --no-autoupdate tunnel --config=${config} run "${tunnel_name}"
