# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Connects and runs the Cloudfalred tunnel
# ==============================================================================
certificate="/root/.cloudflared/cert.pem"
config="/root/.cloudflared/config.yml"
tunnel_name="$(bashio::config 'tunnel_name')"

bashio::log.info "Starting Cloudflared Tunnel..."

exec /opt/cloudflared --origincert=${certificate} tunnel --config=${config} run "${tunnel_name}"
