# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Runs the Cloudfalred tunnel for HomeAssistant
# ==============================================================================
declare default_config=/tmp/config.json

if bashio::config.true 'quick_tunnel'; then
    bashio::log.info "Connecting Cloudflared Quick Tunnel..."
    exec cloudflared tunnel --url "http://homeassistant:$(bashio::core.port)"
else

    data_path="/data"
    if bashio::config.has_value 'data_path'; then
        data_path=$(bashio::config 'data_path' | sed 's:/*$::')
    fi
    # Check status file which exists when custom_config exists
    if bashio::fs.file_exists "/tmp/custom_config"; then
        config="${data_path}/config.yml"
        rm /tmp/custom_config
    else
        config="${default_config}"
    fi
    certificate="${data_path}/cert.pem"
    tunnel_name="$(bashio::config 'tunnel_name')"

    bashio::log.info "Connecting Cloudflared Tunnel..."
    bashio::log.debug "using ${config} config file"
    exec cloudflared --origincert="${certificate}" --no-autoupdate tunnel --config="${config}" run "${tunnel_name}"
fi
