# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/init-cloudflared-config/up
/bin/bash -c /etc/s6-overlay/scripts/10-cloudflared-config.sh