# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/init-cloudflared-log-level/up
/bin/bash -c /etc/s6-overlay/scripts/cloudflared-log-level.sh