# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Caddy..."
exec env XDG_DATA_HOME=/data caddy run --config /etc/caddy/Caddyfile
