# cloudflared/rootfs/etc/services.d/npm_connect/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Connects and runs the Cloudfalred tunnel to connect to Nginxproxymanager
# ==============================================================================
certificate="/data/cert.pem"
config_npm="/data/config_npm.yml"
tunnel_name="$(bashio::config 'tunnel_name')"

bashio::log.info "Starting Cloudflared Tunnel..."

if bashio::config.true 'nginxproxymanager' ; then
    exec /opt/cloudflared --origincert=${certificate} tunnel --config=${config_npm} run "${tunnel_name}"
fi

bashio::log.info "No NPM"
