# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/cloudflared/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflared
# Runs the Cloudflare Tunnel for Home Assistant
# ==============================================================================
declare default_config=/tmp/config.json

if bashio::config.has_value 'tunnel_token' ; then
    tunnel_token=$(bashio::config 'tunnel_token')
    if bashio::config.true 'post_quantum' ; then
        exec cloudflared --no-autoupdate --post-quantum \
            tunnel --metrics="0.0.0.0:36500" --loglevel="${CLOUDFLARED_LOG}" run --token="${tunnel_token}"
    else
        exec cloudflared --no-autoupdate \
            tunnel --metrics="0.0.0.0:36500" --loglevel="${CLOUDFLARED_LOG}" run --token="${tunnel_token}"
    fi
else
    data_path="/data"
    config="${default_config}"
    
    certificate="${data_path}/cert.pem"
    tunnel_name="$(bashio::config 'tunnel_name')"

    bashio::log.info "Connecting Cloudflare Tunnel..."
    bashio::log.debug "using ${config} config file"

    if bashio::config.true 'post_quantum' ; then
        exec cloudflared --origincert=${certificate} --no-autoupdate --post-quantum \
            tunnel --config=${config} --metrics="0.0.0.0:36500" --loglevel="${CLOUDFLARED_LOG}" run "${tunnel_name}"
    else
        exec cloudflared --origincert=${certificate} --no-autoupdate \
            tunnel --config=${config} --metrics="0.0.0.0:36500" --loglevel="${CLOUDFLARED_LOG}" run "${tunnel_name}"
    fi
fi
