# cloudflared/rootfs/etc/s6-overlay/s6-rc.d/init-log-level/up
/bin/bash -c /etc/s6-overlay/scripts/01-log-level.sh